# bedrock.service
[Unit]
Description=Bedrock Database Server with Core Plugin
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/app/Bedrock
ExecStart=/app/Bedrock/bedrock -fork -plugins Core -pluginDir /app/server/core/lib
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment
Environment=LD_LIBRARY_PATH=/app/server/core/lib

[Install]
WantedBy=multi-user.target
